# systemd/t9numpad.service
[Unit]
Description=T9 numpad predictive input daemon
Documentation=man:t9numpad(1)
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=simple
# Run as root by default so we can grab /dev/input/event* without uaccess rules.
# For a user instance drop UID/GID below and ensure the 'input' group membership.
# User=t9numpad
# Group=input
ExecStart=/usr/bin/t9numpad --config /etc/t9numpad/t9numpad.toml
Restart=on-failure
RestartSec=3s

# Hardening
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/t9numpad
ProtectHome=true
NoNewPrivileges=true
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
